# services/wifi-network-init.service
[Unit]
Description=Apply Wi-Fi configuration from boot network-config
After=local-fs.target systemd-remount-fs.service NetworkManager.service
Wants=network-pre.target NetworkManager.service
Before=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/denodo-pi/tools/wifi_network_init.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
